# workstation/systemd/lmc-ai-update.timer
[Unit]
Description=Periodic signed LMC AI Workstation application update check

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min
RandomizedDelaySec=5min
Persistent=true
Unit=lmc-ai-update.service

[Install]
WantedBy=timers.target
