# workstation/systemd/lmc-ai-privileged.service
[Unit]
Description=LMC AI Workstation privileged allowlist helper
After=local-fs.target
Before=lmc-ai-manager.service lmc-ai-node.service lmc-ai-update.service lmc-ai-rollback.service lmc-ai-full-health.service

[Service]
Type=simple
User=root
Group=lmc-ai
WorkingDirectory=/opt/lmc-ai-workstation/current
Environment=PYTHONPATH=/opt/lmc-ai-workstation/current
ExecStart=/usr/bin/python3 -m workstation.privileged_helper.server
Restart=on-failure
RestartSec=3
RuntimeDirectory=lmc-ai-workstation
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=yes
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/lmc-ai-workstation /etc/lmc-ai-workstation /opt/lmc-ai-workstation /var/lib/lmc-ai-workstation
RestrictAddressFamilies=AF_UNIX
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
